# .config/systemd/user/backup.timer
[Unit]
Description=Create backup

[Timer]
OnCalendar=*-*-* 23:50:00
Persistent=true
Unit=backup.service
 
[Install]
WantedBy=timers.target
