# .config/systemd/user/transmission.timer
[Unit]
Description=transmission-daemon start timer

[Timer]
OnBootSec=120
Unit=transmission.service

[Install]
WantedBy=timers.target
