# .config/systemd/user/update.timer
[Unit]
Description=System update timer

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Unit=update.service

[Install]
WantedBy=timers.target
