# .config/systemd/user/update.timer
[Unit]
Description=System update timer

[Timer]
OnBootSec=600
OnUnitActiveSec=600
Unit=update.service

[Install]
WantedBy=timers.target
