# .config/systemd/user/uptime.timer
[Unit]
Description=uptime timer

[Timer]
OnBootSec=30
OnUnitActiveSec=60
Unit=uptime.service

[Install]
WantedBy=timers.target
