# .config/systemd/user/volume_level.socket
[Unit]
Description=Volume socket for Volume service
[Socket]
ListenStream=/run/user/1000/pulse/cli
Accept=true
[Install]
WantedBy=sockets.target
