# .config/systemd/user/volume_level.service
[Unit]
Description=Volume service
[Service]
ExecStart=-/home/haron/.local/bin/XMVol2
StandardInput=socket
