# .config/systemd/user/wi.socket
[Unit]
Description=Weather_icon Socket
[Socket]
ListenFIFO=/tmp/haron/wi_socket
[Install]
WantedBy=sockets.target
