# .config/systemd/user/wi.service
[Unit]
Description=Weather_icon service
[Service]
ExecStart=-/home/haron/.local/bin/XMWI
StandardInput=socket
