# etc/systemd/system/updatedb.timer.d/updatedb.timer.conf
[Timer]
OnCalender=
OnCalendar=21:05
OnBootSec=660
