# .local/share/user/zshist.timer
[Unit]
Description=zshistory synchrinize timer

[Timer]
OnBootSec=1
OnUnitActiveSec=59
Unit=zshist.service

[Install]
WantedBy=timers.target
