# etc/systemd/system/man-db.timer.d/man-db.timer.conf
[Timer]
OnCalender=
OnCalendar=21:00
OnBootSec=600
