# etc/systemd/system/sysctl-reload.timer
[Unit]
Description=sysctl reload timer

[Timer]
OnBootSec=10
Unit=sysctl-reload.service

[Install]
WantedBy=timers.target
