# etc/systemd/system/vnstat-reload.timer
[Unit]
Description=vnstat reload timer

[Timer]
OnBootSec=1
OnUnitActiveSec=59
Unit=vnstat-reload.service

[Install]
WantedBy=timers.target
